# root/etc/services.d/openvpn-privado/run
#!/usr/bin/with-contenv bash

/usr/sbin/openvpn \
	--config "/defaults/${PRIVADOENDPOINT}.ovpn" \
	--auth-user-pass /auth.txt